# PiServer/boot/pi_server.service
[Unit]
Description=PiServer web control backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/PiDrive/PiServer
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/pi/PiDrive/PiServer/server.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
